# .config/systemd/user/speaker-agent.service
[Unit]
Description=Bluetooth speaker agent

[Service]
ExecStart=python speaker-agent.py

[Install]
WantedBy=default.target
